# supervisord.conf
[supervisord]
nodaemon=true

[program:nph-proxy]
command=./nph-proxy.cgi start-fcgi
autostart=true
autorestart=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root

